# SpyderQ_Scripts/SpyderQ_Scripts_backup_20250905/SpyderQ_Scripts/SpyderQ72_IBWatchdog.service
# ===============================================================================
# SPYDER - Autonomous Options Trading System
#
# Spyder Version: 1.0
# Module: SpyderQ72_IBWatchdog.service
# Group: Q (Scripts/Services)
# Purpose: Systemd service for IB Gateway monitoring and auto-restart
# Author: Mohamed Talib
# Date Created: 2025-01-11
# Last Updated: 2025-01-11 Time: 16:20:00
#
# Description:
#     Dedicated systemd service for monitoring IB Gateway health, handling
#     automatic restarts on crashes, managing login timeouts, and ensuring
#     continuous connectivity for all 9 client connections. Works with
#     SpyderB14_MultiClientWatchdog for comprehensive monitoring.
# ===============================================================================

[Unit]
Description=SPYDER IB Gateway Watchdog Service
Documentation=file:///home/adam/Projects/Spyder/docs/ib_watchdog.md
After=network-online.target
Wants=network-online.target
Before=SpyderQ74_SpyderMain.service

# Dependencies
Requires=network-online.target

[Service]
Type=simple
User=adam
Group=adam
WorkingDirectory=/home/adam/Projects/Spyder

# Environment setup
Environment="PYTHONPATH=/home/adam/Projects/Spyder"
Environment="SPYDER_HOME=/home/adam/Projects/Spyder"
Environment="PATH=/home/adam/Projects/Spyder/spyder_venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="IB_GATEWAY_PATH=/home/adam/Jts"
Environment="IB_GATEWAY_PORT_PAPER=4002"
Environment="IB_GATEWAY_PORT_LIVE=4001"
Environment="IB_GATEWAY_RESTART_DELAY=30"

# Virtual environment activation and execution
ExecStartPre=/bin/bash -c 'test -f /home/adam/Projects/Spyder/.env || echo "Warning: .env file not found"'
ExecStartPre=/bin/bash -c 'source /home/adam/Projects/Spyder/spyder_venv/bin/activate'

# Main IB watchdog execution
ExecStart=/home/adam/Projects/Spyder/spyder_venv/bin/python -m SpyderB_Broker.SpyderB14_MultiClientWatchdog

# Gateway health checks
ExecStartPost=/bin/sleep 10
ExecStartPost=/bin/bash -c 'nc -z localhost 4002 || nc -z localhost 4001 || (echo "IB Gateway not responding" && exit 1)'

# Restart configuration
Restart=always
RestartSec=30
StartLimitInterval=600
StartLimitBurst=3

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=60
TimeoutStartSec=120

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryLimit=1G
CPUQuota=75%

# Security settings (less restrictive for IB Gateway control)
PrivateTmp=false
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false
ReadWritePaths=/home/adam/Projects/Spyder /home/adam/Jts

# Logging
StandardOutput=append:/home/adam/Projects/Spyder/logs/system/ib_watchdog.log
StandardError=append:/home/adam/Projects/Spyder/logs/system/ib_watchdog_error.log
SyslogIdentifier=spyder-ib-watchdog

# Success conditions
SuccessExitStatus=0 143

# Special handling for IB Gateway
RestartPreventExitStatus=0

[Install]
WantedBy=multi-user.target
Alias=spyder-ib-watchdog.service