# SpyderQ_Scripts/SpyderQ_Scripts_backup_20250905 (2)/SpyderQ_Scripts/SpyderQ70_Watchdog.service
# ===============================================================================
# SPYDER - Autonomous Options Trading System
#
# Spyder Version: 1.0
# Module: SpyderQ70_Watchdog.service
# Group: Q (Scripts/Services)
# Purpose: Systemd service for Spyder Production Watchdog
# Author: Mohamed Talib
# Date Created: 2025-01-11
# Last Updated: 2025-01-11 Time: 16:00:00
#
# Description:
#     Systemd service unit for the Spyder Production Watchdog that monitors
#     all system components, handles automatic restarts, checks IB Gateway
#     connectivity, and ensures system health. This service starts automatically
#     on boot and restarts components as needed.
# ===============================================================================

[Unit]
Description=SPYDER Production Watchdog Service
Documentation=file:///home/adam/Projects/Spyder/docs/watchdog.md
After=network-online.target
Wants=network-online.target
Before=SpyderQ74_SpyderMain.service

# Dependencies
Requires=network-online.target

[Service]
Type=simple
User=adam
Group=adam
WorkingDirectory=/home/adam/Projects/Spyder

# Environment setup
Environment="PYTHONPATH=/home/adam/Projects/Spyder"
Environment="SPYDER_HOME=/home/adam/Projects/Spyder"
Environment="PATH=/home/adam/Projects/Spyder/spyder_venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"

# Virtual environment activation and execution
ExecStartPre=/bin/bash -c 'test -f /home/adam/Projects/Spyder/.env || echo "Warning: .env file not found"'
ExecStartPre=/bin/bash -c 'source /home/adam/Projects/Spyder/spyder_venv/bin/activate'

# Main watchdog execution
ExecStart=/home/adam/Projects/Spyder/spyder_venv/bin/python /home/adam/Projects/Spyder/SpyderQ_Scripts/SpyderQ24_ProductionWatchdog.py

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
SendSIGKILL=yes

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryLimit=1G
CPUQuota=50%

# Security settings
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/adam/Projects/Spyder/logs /home/adam/Projects/Spyder/pids /home/adam/Projects/Spyder/data

# Logging
StandardOutput=append:/home/adam/Projects/Spyder/logs/system/watchdog.log
StandardError=append:/home/adam/Projects/Spyder/logs/system/watchdog_error.log
SyslogIdentifier=spyder-watchdog

# Health check
ExecReload=/bin/kill -HUP $MAINPID
SuccessExitStatus=0 143

[Install]
WantedBy=multi-user.target
Alias=spyder-watchdog.service