# Spyder/SpyderQ_Scripts/SpyderQ74_SpyderMain.service
# ===============================================================================
# SPYDER - Autonomous Options Trading System
#
# Service: SpyderQ74_SpyderMain.service  
# Group: Q (Scripts/Services)
# Purpose: Main systemd service for complete Spyder Trading System
# Author: Mohamed Talib
# Date Created: 2025-01-11
# Last Updated: 2025-01-11 Time: 12:00:00
#
# Description:
#     Primary systemd service that manages the entire Spyder trading system
#     including IB Gateway monitoring, multi-client connections, AI agents,
#     and dashboard components. Ensures proper startup order and dependencies.
# ===============================================================================

[Unit]
Description=SPYDER Autonomous Options Trading System - Main Service
Documentation=file:///home/adam/Projects/Spyder/docs/README.md
After=network-online.target
Wants=network-online.target
Requires=SpyderQ70_Watchdog.service SpyderQ71_Metrics.service

[Service]
Type=simple
User=adam
Group=adam
WorkingDirectory=/home/adam/Projects/Spyder

# Load secrets from .env (API keys, tokens, mode flags)
EnvironmentFile=/home/adam/Projects/Spyder/.env

# Environment variables
Environment="PYTHONPATH=/home/adam/Projects/Spyder"
Environment="SPYDER_HOME=/home/adam/Projects/Spyder"
Environment="PATH=/home/adam/Projects/Spyder/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="QT_QPA_PLATFORM=offscreen"
Environment="PYTHONUNBUFFERED=1"

# Pre-start checks
ExecStartPre=/bin/bash -c 'test -f /home/adam/Projects/Spyder/.env || (echo "ERROR: .env file missing" && exit 1)'
ExecStartPre=/bin/bash -c '/home/adam/Projects/Spyder/.venv/bin/python -c "import SpyderA_Core" || (echo "ERROR: Python environment invalid" && exit 1)'

# Main execution - Start all components
ExecStart=/home/adam/Projects/Spyder/.venv/bin/python -m SpyderA_Core.SpyderA06_MasterController

# Post-start verification
ExecStartPost=/bin/sleep 5
ExecStartPost=/bin/bash -c 'curl -sf localhost:8000/metrics > /dev/null || echo "Warning: Prometheus metrics not yet available"'

# Stop command
ExecStop=/bin/kill -SIGTERM $MAINPID

# Restart behavior
Restart=on-failure
RestartSec=30
StartLimitInterval=300
StartLimitBurst=5

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryLimit=4G
CPUQuota=200%

# Security
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/adam/Projects/Spyder/logs /home/adam/Projects/Spyder/data

# Logging
StandardOutput=append:/home/adam/Projects/Spyder/logs/system/spyder_main.log
StandardError=append:/home/adam/Projects/Spyder/logs/system/spyder_main_error.log

# Timeouts
TimeoutStartSec=120
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
Alias=spyder.service