# Spyder/SpyderQ_Scripts/SpyderQ75_SessionStart.timer
# ===============================================================================
# SPYDER - Autonomous Options Trading System
#
# Unit: SpyderQ75_SessionStart.timer
# Group: Q (Scripts/Services)
# Purpose: Ensure the Spyder main service is running before market open
#          Monday–Friday at 08:55 ET (host-layer redundancy atop the in-app A04 gate)
# ===============================================================================

[Unit]
Description=SPYDER Session Start Timer — weekday 08:55 ET
Documentation=file:///home/adam/Projects/Spyder/docs/README.md

[Timer]
# Fire weekdays at 08:55 America/New_York
OnCalendar=Mon..Fri *-*-* 08:55:00 America/New_York
Persistent=false
Unit=SpyderQ75_SessionStart.service

[Install]
WantedBy=timers.target
