# Spyder/SpyderQ_Scripts/SpyderQ71_Metrics.service
# ===============================================================================
# SPYDER - Autonomous Options Trading System
#
# Spyder Version: 1.0
# Module: SpyderQ71_Metrics.service
# Group: Q (Scripts/Services)
# Purpose: Systemd service for Prometheus Metrics Collector
# Author: Mohamed Talib
# Date Created: 2025-01-11
# Last Updated: 2025-01-11 Time: 16:10:00
#
# Description:
#     Systemd service unit for the Prometheus Metrics Collector that monitors
#     all 9 IB Gateway client connections, system performance, trading metrics,
#     and exposes data on port 9090 for dashboard integration. Works with
#     SpyderB15_PrometheusMetrics and SpyderG07_PrometheusMetricsDisplay.
# ===============================================================================

[Unit]
Description=SPYDER Prometheus Metrics Collector Service
Documentation=file:///home/adam/Projects/Spyder/docs/metrics.md
After=network-online.target SpyderQ70_Watchdog.service
Wants=network-online.target
Before=SpyderQ74_SpyderMain.service

# Dependencies
Requires=network-online.target

[Service]
Type=simple
User=adam
Group=adam
WorkingDirectory=/home/adam/Projects/Spyder

# Load secrets from .env
EnvironmentFile=/home/adam/Projects/Spyder/.env

# Environment setup
Environment="PYTHONPATH=/home/adam/Projects/Spyder"
Environment="SPYDER_HOME=/home/adam/Projects/Spyder"
Environment="PATH=/home/adam/Projects/Spyder/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="PROMETHEUS_PORT=9090"
Environment="PROMETHEUS_CLIENT_ID=9"

# Virtual environment activation and execution
ExecStartPre=/bin/bash -c 'test -f /home/adam/Projects/Spyder/.env || echo "Warning: .env file not found"'
ExecStartPre=/bin/bash -c 'test -x /home/adam/Projects/Spyder/.venv/bin/python || (echo "ERROR: .venv not found" && exit 1)'
ExecStartPre=/bin/bash -c 'nc -z localhost 9090 && echo "Warning: Port 9090 already in use" && exit 1 || exit 0'

# Main metrics collector execution
ExecStart=/home/adam/Projects/Spyder/.venv/bin/python -m SpyderB_Broker.SpyderB15_PrometheusMetrics

# Health check endpoint
ExecStartPost=/bin/sleep 5
ExecStartPost=/bin/bash -c 'curl -s http://localhost:9090/metrics > /dev/null || exit 1'

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
TimeoutStartSec=60

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryLimit=2G
CPUQuota=100%

# Network configuration
PrivateNetwork=false
RestrictAddressFamilies=AF_INET AF_INET6

# Security settings
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/adam/Projects/Spyder/logs /home/adam/Projects/Spyder/data /home/adam/Projects/Spyder/metrics

# Logging
StandardOutput=append:/home/adam/Projects/Spyder/logs/system/prometheus_metrics.log
StandardError=append:/home/adam/Projects/Spyder/logs/system/prometheus_metrics_error.log
SyslogIdentifier=spyder-metrics

# Success conditions
SuccessExitStatus=0 143

[Install]
WantedBy=multi-user.target
Alias=spyder-metrics.service