# Spyder/SpyderQ_Scripts/SpyderQ78_ZeroDteGuard.service
# ===============================================================================
# SPYDER - Autonomous Options Trading System
#
# Service: SpyderQ78_ZeroDteGuard.service
# Group: Q (Scripts/Services)
# Purpose: Host-layer enforcement of the Zero-DTE no-new-risk cutoff at 15:45 ET.
#          Emits a FLATTEN_REQUEST event via the SpyderA04 scheduler CLI helper
#          to reinforce the in-app cutoff gate.
#          Triggered by SpyderQ78_ZeroDteGuard.timer at 15:45 ET weekdays.
# ===============================================================================

[Unit]
Description=SPYDER Zero-DTE No-New-Risk Cutoff Guard
Documentation=file:///home/adam/Projects/Spyder/docs/README.md

[Service]
Type=oneshot
User=adam
Group=adam
WorkingDirectory=/home/adam/Projects/Spyder
EnvironmentFile=/home/adam/Projects/Spyder/.env
Environment="PYTHONPATH=/home/adam/Projects/Spyder"
Environment="PATH=/home/adam/Projects/Spyder/.venv/bin:/usr/local/bin:/usr/bin:/bin"

# Emit a zero_dte_cutoff signal to the running process via a sentinel file
# that SpyderD31 monitors.  A missing process is silently ignored.
ExecStart=/bin/bash -c '\
  LOCK_DIR=/tmp/spyder_signals && \
  mkdir -p "$LOCK_DIR" && \
  echo "{\"type\":\"zero_dte_no_new_risk\",\"ts\":\"$(date -u +%%Y-%%m-%%dT%%H:%%M:%%SZ)\"}" \
    > "$LOCK_DIR/zero_dte_cutoff.json" && \
  echo "Zero-DTE cutoff signal written to $LOCK_DIR/zero_dte_cutoff.json"'
